# backend/release/canary.service
[Unit]
Description=The canary backend server
After=network.target

[Service]
ExecStart=/usr/sbin/canary -c /usr/local/etc/canary/canary.json
StandardOutput=append:/var/log/canary/access.log
StandardError=append:/var/log/canary/error.log
ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGQUIT
Restart=always
RestartSec=5
TimeoutStopSec=5
KillMode=mixed
PrivateTmp=true

[Install]
WantedBy=multi-user.target
