# backend/release/container-swiftfun.service
# container-swiftfun.service
# autogenerated by Podman 3.2.3
# Mon Nov 22 11:27:22 CST 2021

[Unit]
Description=Podman container-swiftfun.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start swiftfun
ExecStop=/usr/bin/podman stop
Type=forking

[Install]
WantedBy=multi-user.target default.target